# content/etc/cron.d/autopull
* * * * * root /usr/local/bin/autopull.sh > /proc/1/fd/1 2>&1
