# content/etc/cron.d/autopull
* * * * * root /usr/local/bin/autopull.sh
