# content/etc/cron.d/tasmota
*/5 * * * * root /usr/local/munin/tasmota.pl cache
