# content/etc/service/smokeping/finish
#!/bin/sh

echo "Smokeping service stopped at $(date)"
