# content/etc/service/smokeping/finish
#!/bin/sh

ps faux | grep usr.sbin.smokeping | awk '{print $2}' | xargs kill
