# content/etc/service/smokeping/finish
#!/bin/sh

ps faux | grep -- "[s]mokeping --config" | awk '{print $2}' | xargs kill
