# content/etc/service/smokeping/run
#!/bin/sh

exec /usr/bin/perl /usr/sbin/smokeping --config=/etc/smokeping/config --nodaemon --logfile=/var/log/smokeping.log
