# content/etc/service/icinga/finish
#!/bin/sh

echo "Icinga service stopped at $(date)"
