# content/etc/service/icinga/run
#!/bin/sh

/usr/sbin/icinga /etc/icinga/icinga.cfg
