# content/etc/service/icinga/run
#!/bin/sh

exec /usr/sbin/icinga /etc/icinga/icinga.cfg
