# content/etc/service/php-fpm/finish
#!/bin/sh

echo "Php-fpm service stopped at $(date)"
