# content/etc/service/munin-node/finish
#!/bin/sh

ps faux | grep perl.*munin-nod[e] | awk '{print $2}' | xargs kill
