# content/etc/service/munin-node/finish
#!/bin/sh

echo "Munin-node service stopped at $(date)"
