# content/etc/service/munin-node/run
#!/bin/sh

exec /usr/bin/perl -wT /usr/sbin/munin-node
