# content/etc/service/munin-node/run
#!/bin/sh

/usr/bin/perl -wT /usr/sbin/munin-node
