# content/etc/service/fcgiwrap/finish
#!/bin/sh

echo "Fcgiwrap service stopped at $(date)"
