# content/etc/service/fcgi-smokeping/finish
#!/bin/sh

echo "Fcgi-smokeping service stopped at $(date)"
